# dev_build/client_condor/root/etc/sv/autofs/run
#!/bin/bash

BROWSE_MODE="no" /usr/sbin/automount -t 0 -f /etc/auto.master
