# dev_build/server/root/etc/sv/pyglidein_server/run
#!/bin/bash
exec 2>&1
exec chpst -u pyglidein /etc/sv/pyglidein_server/start_pyglidein_server.sh
