# dev_build/server/root/etc/sv/condor/run
#!/bin/bash
exec 2>&1
exec chpst -u condor /etc/sv/condor/start_condor.sh
